# dist/upstart/pantry-client.conf
# /etc/init/pantry-client.conf
# Upstart script for the Pantry Client

description "Pantry Client"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 5 10

exec pantry-client -c /etc/pantry/client.yml
